# supervisord.conf
[unix_http_server]          ; supervisord的unix socket服务配置
file=~/supervisor_kst.sock   ; socket文件的保存目录

[supervisord]
logfile=~/supervisord_kst.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=~/supervisord_kst.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false             ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://~/supervisor_kst.sock ; use a unix:// URL  for a unix socket

[program:kst_server]
command=python kst_server.py
autorestart=true            ; retstart at unexpected quit (default: true)
startretries=10   
stopwaitsecs=1             ; max num secs to wait before SIGKILL (default 10)
stopsignal=TERM
log_stdout=true       
log_stderr=true        
logfile=~/kst_log  
logfile_maxbytes=50MB    


